# siren.service
# systemd 서비스 유닛 — 세이렌 크롤링 대시보드
#
# 설치(서버에서):
#   sudo cp siren.service /etc/systemd/system/siren.service
#   (경로·사용자명이 다르면 아래 값을 수정)
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now siren
#   sudo systemctl status siren        # 상태 확인
#   journalctl -u siren -f             # 로그 실시간 확인
#
# 코드 갱신 후 재시작:
#   git pull && sudo systemctl restart siren

[Unit]
Description=Siren Crawling Dashboard
After=network.target

[Service]
Type=simple
# 서버의 실제 사용자명으로 변경 (Ubuntu 이미지는 보통 ubuntu, Oracle Linux 는 opc)
User=ubuntu
# 프로젝트를 클론한 실제 절대경로로 변경 (저장소명이 S_Crawling 이면 아래대로)
WorkingDirectory=/home/ubuntu/S_Crawling
# 가상환경 안의 gunicorn 절대경로. --workers 1 필수(스케줄러 중복 방지)
ExecStart=/home/ubuntu/S_Crawling/.venv/bin/gunicorn --workers 1 --threads 4 --timeout 120 --bind 0.0.0.0:5000 wsgi:app
Restart=always
RestartSec=5
# 서비스키를 환경변수로 줄 경우(service_key.txt 파일을 쓰면 생략 가능)
# Environment=DATA_GO_KR_SERVICE_KEY=여기에키

[Install]
WantedBy=multi-user.target
